# ops/systemd/cloudflared-providers.service
[Unit]
Description=cloudflared tunnel for PolyStore provider hostnames
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared --no-autoupdate --config %h/.config/cloudflared/providers.polynomialstore.com.yml tunnel run
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
