# ops/systemd/nilchaind.service
[Unit]
Description=NilStore devnet chain daemon (nilchaind)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/nilstore/nilchaind.env
WorkingDirectory=/opt/nilstore
ExecStart=/bin/bash -lc 'exec "$NILCHAIND_BIN" start --home "$NIL_HOME" --rpc.laddr "$NIL_RPC_LADDR" --minimum-gas-prices "$NIL_GAS_PRICES" --api.enable'
Restart=on-failure
RestartSec=2
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
