# ops/systemd/polystore-gateway-router.service
[Unit]
Description=PolyStore router gateway (polystore_gateway)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/polystore/polystore-gateway-router.env
WorkingDirectory=/opt/polystore
ExecStart=/bin/bash -lc 'exec "$NIL_GATEWAY_BIN"'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
