# ops/systemd/polystore-faucet.service
[Unit]
Description=PolyStore devnet faucet (polystore_faucet)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/polystore/polystore-faucet.env
WorkingDirectory=/opt/polystore
ExecStart=/bin/bash -lc 'exec "$NIL_FAUCET_BIN"'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
